# systemd/housing-daily.timer
[Unit]
Description=Prepare housing notices weekdays before the morning briefing

[Timer]
OnCalendar=Mon..Fri *-*-* 06:30:00 Asia/Seoul
Persistent=false
RandomizedDelaySec=0
Unit=housing-daily.service

[Install]
WantedBy=timers.target
