# systemd/housing-daily.service
[Unit]
Description=Daily Seoul housing notice crawl and Telegram report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/data/crawler
EnvironmentFile=/home/ubuntu/.config/monitor-platform/telegram.env
ExecStart=/usr/local/bin/node --no-network-family-autoselection /data/crawler/src/daily.mjs
TimeoutStartSec=2h
NoNewPrivileges=true
PrivateTmp=true
