# systemd/housing-daily.service
[Unit]
Description=Prepare Seoul housing notices and AI reviews
After=network-online.target
Wants=network-online.target
OnFailure=monitor-failure-notify@%n.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/data/crawler
EnvironmentFile=/home/ubuntu/.config/monitor-platform/telegram.env
EnvironmentFile=-/home/ubuntu/.config/monitor-platform/openai-api.env
Environment=CODEX_API_FALLBACK_ENABLED=true
Environment=CODEX_API_DAILY_CALL_LIMIT=10
Environment=DAILY_REPORT_ENABLED=false
ExecStart=/usr/bin/flock --wait 7200 /data/crawler/data/batch.lock /usr/local/bin/node --no-network-family-autoselection /data/crawler/src/daily.mjs
TimeoutStartSec=2h
NoNewPrivileges=true
PrivateTmp=true
UMask=0077
