# systemd/jobs-daily.service
[Unit]
Description=Daily job notice collection, filtering, and Telegram digest
After=network-online.target
Wants=network-online.target
OnFailure=monitor-failure-notify@%n.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/data/crawler
EnvironmentFile=/home/ubuntu/.config/monitor-platform/telegram.env
ExecStart=/usr/bin/xvfb-run -a /usr/local/bin/node --no-network-family-autoselection /data/crawler/src/jobs-daily.mjs
TimeoutStartSec=2h
NoNewPrivileges=true
PrivateTmp=true
