# systemd/jobs-daily.service
[Unit]
Description=Daily job notice collection, filtering, and Telegram digest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/data/crawler
EnvironmentFile=/home/ubuntu/.config/monitor-platform/telegram.env
ExecStart=/usr/bin/xvfb-run -a /usr/local/bin/node /data/crawler/src/jobs-daily.mjs
TimeoutStartSec=2h
NoNewPrivileges=true
PrivateTmp=true
