# systemd/jobs-daily.timer
[Unit]
Description=Prepare job notices weekdays before the morning briefing

[Timer]
OnCalendar=Mon..Fri *-*-* 07:40:00 Asia/Seoul
Persistent=false
RandomizedDelaySec=0
Unit=jobs-daily.service

[Install]
WantedBy=timers.target
