# systemd/monitor-backup.service
[Unit]
Description=Back up Life Daemon SQLite state
After=network-online.target
OnFailure=monitor-failure-notify@%n.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/data/crawler
EnvironmentFile=/home/ubuntu/.config/monitor-platform/telegram.env
ExecStart=/usr/local/bin/node /data/crawler/src/admin/backup-data.mjs
TimeoutStartSec=30min
NoNewPrivileges=true
PrivateTmp=true
UMask=0077
