# systemd/morning-briefing.service
[Unit]
Description=Send one combined Life Daemon morning briefing
After=network-online.target
Wants=network-online.target
OnFailure=monitor-failure-notify@%n.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/data/crawler
EnvironmentFile=/home/ubuntu/.config/monitor-platform/telegram.env
ExecStart=/usr/local/bin/node --no-network-family-autoselection /data/crawler/src/morning-briefing.mjs
TimeoutStartSec=65min
NoNewPrivileges=true
PrivateTmp=true
UMask=0077
